# deployment/performance/supervisor-queue.conf
[program:license-queue]
command=php /var/www/obd2sw/backend/artisan queue:work --queue=default --tries=3 --sleep=3 --timeout=120
directory=/var/www/obd2sw/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/license-queue.log
